# data/scripts/e2guardian.service.in
[Unit]
Description=E2guardian Web filtering
After=network.target

[Service]
Type=forking
IgnoreSIGPIPE=no
GuessMainPID=no
ExecStart=@DGBINDIR@/e2guardian
ExecReload=@DGBINDIR@/e2guardian -r
ExecStop=@DGBINDIR@/e2guardian -q

[Install]
WantedBy=multi-user.target
