# data/scripts/e2guardian.service.in
[Unit]
Description=E2guardian Web filtering
After=network.target
TimeoutStartSec=20s

[Service]
LimitNOFILE=65535
Type=forking
DefaultTasksMax=infinity
IgnoreSIGPIPE=no
GuessMainPID=no
ExecStart=@DGBINDIR@/e2guardian
ExecReload=@DGBINDIR@/e2guardian -r
UMask=027

[Install]
WantedBy=multi-user.target
