# data/scripts/e2guardianrasp.service.in
#  This file is part of systemd.
#  Must be used with rasp2/3 and e2 version > 5.1  

[Unit]
Description=E2guardian Web filtering
After=network.target

[Service]
TasksMax=infinity
StackLIMIT=infinity:infinity
LimitNOFILE=65535
Type=forking
IgnoreSIGPIPE=no
GuessMainPID=no
ExecStart=/usr/sbin/e2guardian
ExecReload=/usr/sbin/e2guardian -r
UMask=027
LimitSTACK=1048576

[Install]
WantedBy=multi-user.target
