# data/scripts/e2guardian.service.in
[Unit]
Description=E2guardian Web filtering
After=network.target

[Service]
Type=forking
ExecStart=@DGBINDIR@/e2guardian
ExecReload=/usr/sbin/e2guardian -r

[Install]
WantedBy=multi-user.target
